# root/load-iptables.service
[Unit]
Description=Load iptables firewall rules

[Service]
Type=oneshot
ExecStart=/root/iptables.rules

[Install]
WantedBy=multi-user.target
